# docker/php/config/laravel.cron
*/1 * * * * php /srv/artisan schedule:run
